# gopls-cpulimit/gopls-cpulimit.service
[Unit]
Description=cpulimit to gopls process

[Service]
Type=simple
ExecStart=/usr/local/bin/gopls-cpulimit.sh

[Install]
WantedBy=multi-user.target
