# domain/update-domain.service
[Unit]
Description=Update dns record

[Service]
Type=simple
ExecStart=/usr/local/bin/update-domain.sh
